# ubus/systemd/ubus.service
[Unit]
Description=OpenWrt micro bus
Requires=ubus.socket

[Service]
ExecStart=/usr/local/sbin/ubusd
