# ubus/systemd/ubus.socket
[Unit]
Description=OpenWrt micro bus socket

[Socket]
ListenStream=/var/run/ubus.sock

[Install]
WantedBy=sockets.target
